# fancontrol.service
[Unit]
Description=Fan Control Web Server
After=network.target

[Service]
Type=simple
User=mario
WorkingDirectory=/home/mario/Code/FanControl
ExecStart=/usr/bin/python3 /home/mario/Code/FanControl/simple_server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
